# schellenberg-usb-hack/rootfs/etc/s6-overlay/s6-rc.d/frontend/run
#!/usr/bin/with-contenv bashio

bashio::log.info "Schellenberg Frontend starting..."

cd /opt/schellenberg-nuxt
exec node ./.output/server/index.mjs